# root/etc/services.d/php-fpm/run
#!/bin/bash

exec /usr/bin/php-fpm -F -y /config/nginx/nginx-fpm.conf
