# root/etc/services.d/php-fpm/run
#!/bin/bash

exec /usr/sbin/php-fpm7 -F -y /config/nginx/nginx-fpm.conf
