# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

if pgrep -f "[n]ginx:" > /dev/null; then
  pkill [n]ginx:
fi

exec /usr/sbin/nginx -c /config/nginx/nginx.conf
